# crontab




PATH=/bin:/sbin:/usr/bin:/usr/sbin:/frankbarton/pipeline/tools/scripts/
PYTHONPATH=/frankbarton/pipeline/tools/python


*/5 * * * *	/scripts/maintainJobPermissions.sh
